# units/nvidia-powersave.service
[Unit]
Description="powersave helper"

[Service]
Type=oneshot
ExecStart=-/sbin/rmmod nvidia_drm nvidia-modeset nvidia_uvm nvidia

