# units/nvidia-powersave.timer
[Unit]
Description="powersave timer"

[Timer]
OnActiveSec=1s
OnUnitActiveSec=2s

[Install]
WantedBy=timers.target
